# src/roles/minecraft_tools/files/mc-monitor.service
[Service]
ExecStart=/usr/local/bin/mc-monitor
Restart=on-failure
User=root
Group=root

[Install]
WantedBy=multi-user.target
